# template/master/override/lib/systemd/system/{{package}}.service
[Unit]
Description=Job that runs the {{package}} daemon

[Service]
Type=simple
ExecStart=/opt/venvs/{{package}}/bin/{{PackageService}}
