# override/lib/systemd/system/subfind.service
[Unit]
Description=Job that runs the subfind daemon

[Service]
Type=simple
ExecStart=/opt/venvs/subfind/bin/subfind-web
